# arangodb/arangodb.initd
#!/sbin/openrc-run

command="/usr/sbin/arangod"
command_args="--configuration=/etc/arangodb3/arangod.conf --pid-file=/run/arangodb3/arangod.pid"
pidfile="/run/arangodb3/arangod.pid"

depend() {
  need net
}
